# .kde4/share/autostart/plasma-netbook.desktop
[Desktop Entry]
Exec=plasma-netbook
Hidden=true
OnlyShowIn=KDE;
X-KDE-autostart-phase=BaseDesktop
